# deb-extracted/DEBIAN/postrm
#!/bin/bash
# postrm — clean up the service on removal.
set -e

TARGET_USER="${SUDO_USER:-}"
if [[ -z "${TARGET_USER}" || "${TARGET_USER}" == "root" ]]; then
    TARGET_USER="$(getent passwd 1000 | cut -d: -f1 || true)"
fi

# Disable and remove the per-user service unit.
if [[ -n "${TARGET_USER}" ]]; then
    TARGET_UID="$(id -u "${TARGET_USER}" 2>/dev/null || true)"
    TARGET_HOME="$(getent passwd "${TARGET_USER}" | cut -d: -f6 || true)"
    if [[ -n "${TARGET_UID}" ]]; then
        sudo -u "${TARGET_USER}" XDG_RUNTIME_DIR="/run/user/${TARGET_UID}" \
            systemctl --user disable lorawan-tracker.service 2>/dev/null || true
    fi
    if [[ "$1" == "purge" && -n "${TARGET_HOME}" ]]; then
        rm -f "${TARGET_HOME}/.config/systemd/user/lorawan-tracker.service"
        rm -rf "${TARGET_HOME}/.config/lorawan-tracker"
    fi
fi

# Remove any legacy system-level service.
if [[ -d /run/systemd/system ]] && command -v systemctl &>/dev/null; then
    systemctl disable lorawan-tracker.service 2>/dev/null || true
    rm -f /etc/systemd/system/lorawan-tracker.service 2>/dev/null || true
    systemctl daemon-reload 2>/dev/null || true
fi

exit 0
